# launchd/com.garmin-health-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.garmin-health-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/might/clawd/repos/garmin-health-sync/run.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/might/clawd/repos/garmin-health-sync</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/might/clawd/repos/garmin-health-sync/sync.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/might/clawd/repos/garmin-health-sync/sync.error.log</string>
</dict>
</plist>
